# roles/_http_logrotate/templates/config/crontab
SHELL='/usr/bin/bash'
BASH_ENV='/etc/environment'

59 23 * * * root /usr/sbin/logrotate {{ http_logrotate_config_path }} &> /dev/null
